# 00_Holding/automation/com.user.sync-cto-brain.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.sync-cto-brain</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/clchok/Automation/sync_cto_brain.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/clchok/Library/Mobile Documents/com~apple~CloudDocs/cto-brain</string>
    </array>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/clchok/Automation/logs/cto_brain_launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/clchok/Automation/logs/cto_brain_launchd_error.log</string>
</dict>
</plist>
